# _ARCHIVE_OLD_FILES/debian/phaze-vpn/opt/phaze-vpn/web-portal/phazevpn-portal.service
[Unit]
Description=PhazeVPN Web Portal (Gunicorn)
After=network.target

[Service]
Type=notify
User=www-data
Group=www-data
WorkingDirectory=/opt/phaze-vpn/web-portal
Environment="PATH=/usr/bin:/usr/local/bin"
Environment="VPN_SERVER_IP=phazevpn.com"
Environment="VPN_SERVER_PORT=1194"
Environment="HTTPS_ENABLED=true"
ExecStart=/usr/local/bin/gunicorn --workers 4 --bind 127.0.0.1:5000 --timeout 120 --access-logfile /var/log/phazevpn-portal-access.log --error-logfile /var/log/phazevpn-portal-error.log app:app
ExecReload=/bin/kill -s HUP $MAINPID
Restart=always
RestartSec=5

[Install]
WantedBy=multi-user.target

