# _ARCHIVE_OLD_FILES/debian/phaze-vpn-download.service
[Unit]
Description=PhazeVPN Client Download Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/phaze-vpn
ExecStart=/usr/bin/python3 /opt/phaze-vpn/client-download-server.py
Restart=always
RestartSec=5
StartLimitInterval=0
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target

