# debian/phaze-vpn.service
[Unit]
Description=PhazeVPN Server
After=network.target

[Service]
Type=simple
User=root
WorkingDirectory=/opt/phaze-vpn
ExecStart=/usr/bin/python3 /opt/phaze-vpn/vpn-manager.py start
ExecStop=/usr/bin/python3 /opt/phaze-vpn/vpn-manager.py stop
Restart=on-failure
RestartSec=5

[Install]
WantedBy=multi-user.target

